# app/deploy/supervisord.conf
[program:Timers]
command: /usr/bin/python -m SimpleHTTPServer
directory: /home/ubuntu/Timers/app
autostart: true
stderr_redirect: true
stdout_logfile: /tmp/TimersSitelog